# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/hotplug.d/nat/10-nat-dnsproxy
#!/bin/sh

. /usr/lib/dnsproxy/dnsproxy_api.sh

if [ "$NAT_ACT" = "on" ]; then
	dnsproxy_start_deamon "none"
fi