# boards/model_brcm_bcm470x/iplatform_branch/stable_rel_20141125_cloud/opensource/mcproxy/filesystem/etc/hotplug.d/iface/mcproxy.hotplug
#!/bin/sh
# Copyright (c) 2014 The Linux Foundation. All rights reserved.

. /lib/functions.sh

/etc/init.d/mcproxy enabled && [ "$ACTION" = "ifup" ] && {
	local upstream_ready=0

	option_cb(){
		[ "$1" != "upstream" ] && return
		[ "$2" != "$IFNAMES" ] && return
		upstream_ready=1
	}

	IFNAMES=`/sbin/uci -q get network.${INTERFACE}.ifname`
	config_load mcproxy
	[ $upstream_ready != "0" ] && /etc/init.d/mcproxy restart
}

