# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/etc/hotplug.d/firewall/30-dos_protection
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

if [ "$ACTION" = "start" ]; then
	echo "reload: dos_protection"
	fw s_add 4 f dos_protection
	fw s_add 4 f forwarding_rule dos_protection
fi
