# boards/model_brcm_bcm470x/build/product_configs/c1200v1_eu/filesystems/etc/hotplug.d/iface/02-wireless-vlan
[ "$INTERFACE" == "lan" -a "$ACTION" == "ifup" ] && {
    /sbin/wifi vlan
}

[ "$INTERFACE" == "lan" -a "$ACTION" == "ifdown" ] && {
    : # Null command
}
