# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/etc/hotplug.d/iface/50-improxy
#!/bin/sh
# Copyright (c) 2014 The Linux Foundation. All rights reserved.

. /lib/functions.sh

/etc/init.d/improxy enabled && [ "$ACTION" = "ifup" ] && {
	/etc/init.d/improxy restart
}

