# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/etc/hotplug.d/iface/70-imb
#!/bin/sh

[ "$DEVICE" = "lo" ] && exit 0

[ "$INTERFACE" = "lan" ] && /etc/init.d/imb restart
