# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/access-control.list
/lib/access_control/core.sh
/lib/access_control/core_global.sh
/lib/access_control/core_log.sh
/lib/upgrade/keep.d/access-control
/etc/init.d/access_control
/lib/access_control/config.sh
/sbin/ac
/etc/hotplug.d/firewall/20-access_control
/etc/config/access_control
