# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/administration.list
/lib/administration/core.sh
/etc/hotplug.d/firewall/50-administration
/etc/init.d/administration
/etc/config/administration
