# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/base-filesystem.list
/etc/openwrt_version
/sbin/sound
/etc/init.d/done
/sbin/network_get_firm
/etc/hotplug2-init.rules
/etc/preinit
/sbin/mount_root
/lib/upgrade/platform.sh
/etc/hotplug2-common.rules
/etc/shells
/etc/protocols
/sbin/loadprofile
/sbin/network_firm
/etc/hotplug.d/button/00-button
/etc/config/system
/etc/sysctl.conf
/sbin/sysupgrade
/lib/upgrade/keep.d/base-filesystem
/sbin/led.sh
/etc/banner
/etc/shadow
/etc/rc.common
/sbin/gratuitous_arp
/etc/init.d/watchdog
/etc/init.d/sysntpd
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/init.d/usb
/etc/services
/sbin/wifi_commit
/sbin/wifi_firm
/etc/init.d/garp
/lib/upgrade/common.sh
/etc/sysupgrade.conf
/etc/uci-defaults/migrate-sysctl
/etc/group
/etc/functions.sh
/sbin/firstboot
/sbin/loadqosdatabase
/bin/login.sh
/sbin/loadconfig
/etc/uci-defaults/migrate-shadow
/etc/hosts
/etc/hotplug.d/usb/10-usb
/lib/functions.sh
/etc/config/network
/sbin/sysdebug
/etc/init.d/rcS
/etc/inittab
/etc/openwrt_release
/sbin/wifi
/etc/profile
/lib/functions/uci-protocol-defaults.sh
/etc/init.d/umount
/lib/functions/service.sh
/lib/upgrade/keep.d/base-files-essential
/etc/passwd
/etc/hotplug.d/ieee1394/10-ieee1394
/sbin/getfirm
/etc/init.d/boot
/lib/functions/gpio-button-hotplug.sh
/lib/functions/boot.sh
/bin/ipcalc.sh
/etc/diag.sh
