# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/client_mgmt.list
/usr/bin/client_mgmt
/etc/init.d/client_mgmt
