# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/dos-protection.list
/lib/dos_protection/core.sh
/etc/init.d/dos_protection
/lib/dos_protection/config.sh
/etc/config/dos_protection
/lib/dos_protection/core_log.sh
/lib/dos_protection/core_global.sh
/etc/hotplug.d/firewall/30-dos_protection
/usr/sbin/dosd
