# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/factory-settings-reset.list
/etc/init.d/factory_settings_reset
/usr/bin/factory_settings_reset
