# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/firewall.list
/lib/firewall/core_interface.sh
/etc/hotplug.d/firewall/10-nat-reflection
/lib/firewall/core_log.sh
/lib/firewall/core_tpcmd.sh
/lib/firewall/core_forwarding.sh
/etc/init.d/firewall
/lib/firewall/config.sh
/etc/config/firewall
/lib/firewall/core.sh
/etc/firewall.user
/lib/firewall/fw.sh
/lib/upgrade/keep.d/firewall
/lib/firewall/core_rule.sh
/lib/firewall/core_redirect.sh
/etc/hotplug.d/iface/20-firewall
/lib/firewall/tpcmd.sh
/lib/firewall/core_init.sh
/sbin/fw
/lib/firewall/uci_firewall.sh
