# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/imb.list
/etc/hotplug.d/iface/70-imb
/etc/config/imb
/etc/init.d/imb
/usr/sbin/imbd
