# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/improxy.list
/etc/config/improxy
/etc/ppp/ip-down.d/99improxy.sh
/etc/init.d/improxy
/usr/sbin/improxy
/etc/ppp/ip-up.d/99improxy.sh
/etc/hotplug.d/iface/50-improxy
/etc/hotplug.d/firewall/50-improxy
