# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/kmod-iptvx.list
/etc/hotplug.d/switch/65-iptv
/etc/hotplug.d/iface/65-iptv
/lib/iptv/iptv_config.lua
/lib/iptv/iptv_network.sh
/etc/hotplug.d/firewall/65-iptv
/etc/init.d/iptv
/lib/modules/iplatform/iptvx.ko
/lib/iptv/iptv_core.sh
/etc/config/iptv
