# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/ledctrl.list
/etc/init.d/ledctrl
/usr/bin/ledcli
/etc/init.d/sys-led-end
/usr/bin/ledctrl
/etc/config/ledctrl
/etc/init.d/sys-led-start
