# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/logmail.list
/etc/init.d/logmail-auto
/usr/sbin/logmail
