# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/luarsa.list
/usr/lib/lua/luarsa.so
/usr/bin/luarsa_keys_gen
/etc/init.d/luarsa_keys_gen
