# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/minidlna.list
/usr/sbin/minidlnad
/etc/init.d/minidlna
/etc/config/minidlna
