# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/nat.list
/lib/nat/nat.sh
/etc/config/nat
/lib/nat/nat_config.sh
/sbin/nat
/etc/hotplug.d/iface/20-nat
/etc/hotplug.d/firewall/25-nat
/lib/nat/nat_vs.sh
/lib/nat/nat_core.sh
/lib/nat/nat_tools.lua
/etc/init.d/nat
/lib/nat/nat_var.sh
/lib/nat/nat_comm.sh
/lib/nat/nat_alg.sh
/lib/nat/nat_pt.sh
/lib/nat/nat_log.sh
/lib/nat/nat_init.sh
/lib/nat/config.sh
/lib/nat/nat_dmz.sh
