# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/netifd.list
/lib/netifd/proto/dhcp.sh
/etc/init.d/network
/lib/netifd/netifd-proto.sh
/etc/init.d/networkx
/etc/hotplug.d/iface/00-netstate
/sbin/ifdown
/lib/netifd/dhcp.script
/usr/share/udhcpc/default.script
/lib/network/config.sh
/etc/hotplug.d/iface/10-sysctl
/sbin/netifd
/sbin/devstatus
/etc/hotplug.d/iface/01-ipv6
/etc/config/portspeed
/sbin/ifstatus
/sbin/ifup
/sbin/portspeed
/etc/hotplug.d/iface/80-tcpmss
/sbin/ifmac
