# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/openvpn-polarssl.list
/etc/config/openvpn
/usr/sbin/checkpsw.sh
/etc/init.d/openvpn
/usr/sbin/openvpn
/lib/upgrade/keep.d/openvpn
