# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/parental-control.list
/etc/tsched.d/parental_control
/lib/parental_control/core_log.sh
/etc/hotplug.d/firewall/21-parental_control
/lib/parental_control/config.sh
/etc/init.d/parental_control
/lib/parental_control/core.sh
/etc/config/parental_control
/lib/upgrade/keep.d/parental-control
/lib/parental_control/core_parental_ctrl.sh
