# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/password-reset.list
/usr/bin/password_reset
/etc/init.d/password_reset
