# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/qos.list
/lib/upgrade/keep.d/qos
/etc/init.d/qos
/lib/qos/config.sh
/lib/qos/core.sh
/etc/config/qos
/lib/qos/core_qos.sh
/lib/qos/core_log.sh
