# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/radvd.list
/usr/sbin/radvd
/etc/init.d/radvd
/etc/config/radvd
