# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/samba36-server.list
/etc/config/samba
/usr/sbin/nmbd
/etc/samba/lognames.deny
/etc/samba/smb.conf.template
/etc/samba/valid.dat
/etc/init.d/samba
/etc/samba/upcase.dat
/etc/samba/smbusers
/usr/sbin/smbpasswd
/usr/sbin/smbd
/usr/sbin/samba_multicall
/etc/samba/lowcase.dat
