# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/switch_led.list
/etc/init.d/switch_led
/usr/sbin/lan_connected
/usr/bin/switch_led
/usr/sbin/wan_connected
