# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/switch-setup.list
/lib/switch/config.sh
/sbin/switch_link_down
/lib/switch/core.sh
/etc/init.d/switch
/etc/init.d/phyport
/sbin/switch_link_up
/etc/config/switch
