# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/system-monitor.list
/usr/sbin/sysmond
/etc/init.d/system-monitor
