# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/time-schedule.list
/usr/sbin/tsched
/usr/sbin/tsched_conf
/etc/init.d/tsched
/usr/lib/lua/tsched_conf.lua
