# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/time-settings.list
/lib/time_settings/set_time
/etc/init.d/time_settings
