# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/tmp-server.list
/usr/bin/tmpServer
/etc/init.d/tmpServer
