# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/traffic-stats.list
/etc/init.d/tfstats
/usr/sbin/tfstats
/etc/config/tfstats
