# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/usbshare.list
/etc/hotplug.d/iface/70-ftpex
/sbin/usbuser
/sbin/ftpex
/etc/config/sharecfg
/sbin/ftp_alias
/etc/init.d/usbshare
/etc/hotplug.d/iface/70-samba
/etc/config/usbshare
/usr/sbin/usbshare
/etc/hotplug.d/iface/70-minidlna
