# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/wide-dhcpv6-server.list
/etc/dhcp6s.conf
/usr/sbin/dhcp6s
/etc/init.d/dhcp6s
/etc/config/dhcp6s
/etc/dhcp6cctlkey
/etc/dhcp6sctlkey
