# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/usr/lib/opkg/info/wireless-button.list
/usr/bin/wireless_button
/etc/init.d/wireless_button
