# boards/model_brcm_bcm470x/build/product_configs/c1200v1_eu/filesystems/etc/rc.d/S50wireless_start
#!/bin/sh /etc/rc.common

START=30

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
     /sbin/wifi init start
     return 0
}
