# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/etc/init.d/client_mgmt
#!/bin/sh /etc/rc.common

SERVICE_DAEMONIZE=1
#SERVICE_WRITE_PID=1

START=70
STOP=70

start()
{
	service_start /usr/bin/client_mgmt
	return 0
}

stop()
{
	service_stop /usr/bin/client_mgmt
	return 0
}
