# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/etc/init.d/luarsa_keys_gen
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=48

start() {
#    luarsa_keys_gen
}

restart() {
#    luarsa_keys_gen
}

reload() {
#    luarsa_keys_gen
}
