# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/etc/init.d/parental_control
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=47

PC_LIBDIR=/lib/parental_control

pc() {
	. $PC_LIBDIR/core.sh
	fw_$1
}

start() {
	pc start
}

stop() {
	pc stop
}

restart() {
	pc restart
}

reload() {
	pc reload
}
