# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/init.d/imb
#!/bin/sh /etc/rc.common

START=66
STOP=66

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
	service_start /usr/sbin/imbd

	return 0
}

stop() {
	service_stop /usr/sbin/imbd

	return 0
}
