# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/init.d/sys-led-end
#!/bin/sh /etc/rc.common

START=99

cycle()
{
	sleep 1
	ledcli STATUS_ON
}

start()
{
	cycle &
}
