# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/init.d/tmpServer
#!/bin/sh /etc/rc.common

START=50

start()
{
	service_start /usr/bin/tmpServer &
}

stop()
{
	service_stop /usr/bin/tmpServer
}
