# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/hotplug.d/iface/25-dnsproxy
#!/bin/sh

. /usr/lib/dnsproxy/dnsproxy_api.sh

if [ "$ACTION" = "ifup" ]; then
	dnsproxy_start_deamon "$INTERFACE"
fi