# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/hotplug.d/iface/70-samba
#!/bin/sh

. /lib/functions.sh

if [ "$INTERFACE" != "lan" ]; then
    exit 0
fi

if [ "$ACTION" = "ifup" ]; then
    config_load usbshare
    config_get samba global samba

    if [ "$samba" == "on" ]; then
        /etc/init.d/samba restart
    fi
fi
