# IPlatform/platform/ibase/luci/src/applications/luci-upnp/root/etc/uci-defaults/luci-upnp
#!/bin/sh
# Copyright (c) 2013 The Linux Foundation. All rights reserved.
# Not a Contribution.

/etc/init.d/miniupnpd enabled

rm -f /tmp/luci-indexcache
exit 0
