# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/hotplug.d/firewall/20-access_control
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

if [ "$ACTION" = "start" ]; then
	/etc/init.d/access_control start
fi
