# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/hotplug.d/firewall/25-nat
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

[ "$ACTION" = "start" ] && {
	local f=/var/nat.switch
    
    [ -e $f ] && rm -f "$f"

    /etc/init.d/nat restart
}
