# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/hotplug.d/firewall/50-administration
#!/bin/sh

. /lib/functions.sh

if [ "$ACTION" = "start" ]; then
	/etc/init.d/administration start
fi
