# IPlatform/platform/staging_dir/target-arm-openwrt-linux-uclibc/root-model_brcm_bcm470x/etc/hotplug.d/switch/65-iptv
#!/bin/sh

[ "$ACTION" = "RESTART" ] && /etc/init.d/iptv restart
