# boards/model_brcm_bcm470x/image/c1200v1_eu/rootfs/etc/rc.d/S30wireless_art
#!/bin/sh /etc/rc.common

START=30

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
     /sbin/wifi_firm
     /sbin/wifi init common
     /sbin/wifi init art
     return 0
}
